# coreos-files/ubuntu_haproxy/files/haproxy.supervisor.conf
[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
autostart=true
autorestart=true
user=root 
